# templates/init_password.service
[Unit]
Description=Init Password for image
After=network.target {{init_service_before}}.service

[Service]
Type=oneshot
TimeoutSec=10s
ExecStart=/usr/bin/bash /tmp/init.sh
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
