# templates/init_password.service
[Unit]
Description=Init Password for image
After=network.target

[Service]
ExecStart=/usr/bin/bash /tmp/init.sh

[Install]
WantedBy=multi-user.target
